# supervisor.ini
[program:pypi]
command=/home/jamiecurle/pypi/venv/bin/pypi-server -p 7001 /home/jamiecurle/pypi/packages
directory=/home/jamiecurle/pypi
user=jamiecurle
autostart=true
autorestart=true
redirect_stderr=true
